# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:webapp]
command=gunicorn -w 4 -k uvicorn.workers.UvicornWorker -b 0.0.0.0:8000 main:web_app
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/webapp.err.log
stdout_logfile=/var/log/webapp.out.log

[program:telegram-bot]
command=python telegram_bot.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/bot.err.log
stdout_logfile=/var/log/bot.out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
